# deploy/systemd/math-for-agents-healthcheck.timer.example
[Unit]
Description=Run math-for-agents release healthcheck every five minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
AccuracySec=30s
Persistent=true
Unit=math-for-agents-healthcheck.service

[Install]
WantedBy=timers.target
